# docker/crontab
# try to grow tiles every minute
* * * * * /var/www/html/docker/requestUrl.sh POST /internal/grow-tiles
